# package/DEBIAN/prerm
#!/bin/bash
set -e
systemctl disable snapstage-firstboot.service || true
exit 0
